# php4pro/supervisor/conf.d/queue.conf
[program:queue]
command=php think queue:listen --queue
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/queue.error.log
stdout_logfile=/var/log/supervisor/queue.log